# etc/moxa/npreal2.rules
SUBSYSTEM=="tty", ACTION=="add", RUN+="/usr/lib/npreal2/driver/mxrmtty"
